# systemd/sp500-sma200-daily-pattern.timer
[Unit]
Description=SP500 and SMA200 daily chart patter

[Timer]
OnCalendar=Mon..Fri 08:30 Europe/Paris

[Install]
WantedBy=timers.target
